# autostart/termlogin-zen-2.6.sh.desktop
[Desktop Entry]
Exec=/home/iskra/termlogin-zen-2.6.sh
Icon=system-run
Name=termlogin-zen-2.6.sh
Path=
Terminal=False
Type=Application
